# data/dbus/couchplay-helper.service
# SPDX-License-Identifier: GPL-3.0-or-later
# SPDX-FileCopyrightText: 2025 CouchPlay Contributors

[Unit]
Description=CouchPlay Privileged Helper Daemon
Documentation=https://github.com/hikaps/couchplay

[Service]
Type=dbus
BusName=io.github.hikaps.CouchPlayHelper
ExecStart=/usr/local/libexec/couchplay-helper
User=root
Group=root

# Security hardening - relaxed for user creation
# NOTE: User creation requires write access to /etc and /home, so we use minimal restrictions
NoNewPrivileges=true
ProtectSystem=false
ProtectHome=false
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=false
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
UMask=0077

# Allow access to input devices
ReadWritePaths=/dev/input
# Allow access to user home directories for PipeWire config
ReadWritePaths=/home
ReadWritePaths=/var/home
# Allow access to /etc for user creation (useradd needs to modify passwd, group, shadow)
ReadWritePaths=/etc

[Install]
WantedBy=multi-user.target
